# steam_killer/steam-killer.service
[Unit]
Description=Steam Killer Daemon
Wants=default.target

[Service]
Type=simple
ExecStart=steam-killer
